# etc/systemd/system/gpg-preset.service
[Unit]
Description=Run gpg-preset script.

[Service]
ExecStart=/root/bin/gpg-preset.zsh
